# data/supervisor.conf
[program:internetradio]
command = gunicorn -b 0.0.0.0:12433 app:app
directory = /opt/loxberry/data/plugins/internetradio/scripts
user = root
environment=HOME="/root", USER="root"
